# console/deploy/linux/s5-console.service
[Unit]
Description=S5 Local Training Console B
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=s5console
WorkingDirectory=/opt/sichuan-mahjong-engine
ExecStart=/opt/sichuan-mahjong-engine/.venv/bin/python -m console.run_server_b --config /etc/s5-console.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
